# config/systemd/cortex-backup.timer
# cortex-backup.timer — daily cortex database backup with randomized delay.
#
# Install alongside cortex-backup.service:
#   systemctl --user enable --now cortex-backup.timer
# Inspect: systemctl --user list-timers cortex-backup.timer

[Unit]
Description=Daily cortex WAL-safe SQLite backup

[Timer]
OnCalendar=daily
# Spread fleet hosts across an hour so they don't all back up at once.
RandomizedDelaySec=1h
# Run after boot if the scheduled time was missed.
Persistent=true

[Install]
WantedBy=timers.target
